# roles/selenium/files/xvbf.upstart
description "xvbf"

start on runlevel [23]
stop on shutdown

env DISPLAY=:99
setuid selenium

exec Xvfb :99 -screen 0 1920x1080x24
